# package/luci-app-lucicodex/root/etc/uci-defaults/40_luci-app-lucicodex
#!/bin/sh

# Initialize lucicodex configuration
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@lucicodex[-1]
	add ucitrack lucicodex
	set ucitrack.@lucicodex[-1].init=lucicodex
	commit ucitrack
EOF

# Create default lucicodex config with named 'main' section if it doesn't exist
if ! uci -q get lucicodex.main >/dev/null 2>&1; then
	# Check for legacy anonymous sections and migrate
	if uci -q get lucicodex.@api[0].key >/dev/null 2>&1; then
		# Migrate from anonymous api section
		OLD_KEY=$(uci -q get lucicodex.@api[0].key)
		OLD_PROVIDER=$(uci -q get lucicodex.@api[0].provider)
		OLD_MODEL=$(uci -q get lucicodex.@api[0].model)
		OLD_ENDPOINT=$(uci -q get lucicodex.@api[0].endpoint)
		OLD_OPENAI=$(uci -q get lucicodex.@api[0].openai_key)
		OLD_ANTHROPIC=$(uci -q get lucicodex.@api[0].anthropic_key)
		OLD_OPENAI_MODEL=$(uci -q get lucicodex.@api[0].openai_model)
		OLD_OPENAI_ENDPOINT=$(uci -q get lucicodex.@api[0].openai_endpoint)
		OLD_ANTHROPIC_MODEL=$(uci -q get lucicodex.@api[0].anthropic_model)
		OLD_ANTHROPIC_ENDPOINT=$(uci -q get lucicodex.@api[0].anthropic_endpoint)
	fi
	
	if uci -q get lucicodex.@settings[0].dry_run >/dev/null 2>&1; then
		# Migrate from anonymous settings section
		OLD_DRY_RUN=$(uci -q get lucicodex.@settings[0].dry_run)
		OLD_CONFIRM=$(uci -q get lucicodex.@settings[0].confirm_each)
		OLD_TIMEOUT=$(uci -q get lucicodex.@settings[0].timeout)
		OLD_MAX_CMD=$(uci -q get lucicodex.@settings[0].max_commands)
		OLD_LOG=$(uci -q get lucicodex.@settings[0].log_file)
	fi

	# Create named 'main' section
	uci -q batch <<-EOF
		set lucicodex.main='settings'
		set lucicodex.main.provider='${OLD_PROVIDER:-gemini}'
		set lucicodex.main.key='${OLD_KEY:-}'
		set lucicodex.main.model='${OLD_MODEL:-gemini-2.5-flash}'
		set lucicodex.main.endpoint='${OLD_ENDPOINT:-https://generativelanguage.googleapis.com/v1beta}'
		set lucicodex.main.openai_key='${OLD_OPENAI:-}'
		set lucicodex.main.openai_model='${OLD_OPENAI_MODEL:-}'
		set lucicodex.main.openai_endpoint='${OLD_OPENAI_ENDPOINT:-}'
		set lucicodex.main.anthropic_key='${OLD_ANTHROPIC:-}'
		set lucicodex.main.anthropic_model='${OLD_ANTHROPIC_MODEL:-}'
		set lucicodex.main.anthropic_endpoint='${OLD_ANTHROPIC_ENDPOINT:-}'
		set lucicodex.main.dry_run='${OLD_DRY_RUN:-1}'
		set lucicodex.main.confirm_each='${OLD_CONFIRM:-0}'
		set lucicodex.main.timeout='${OLD_TIMEOUT:-30}'
		set lucicodex.main.max_commands='${OLD_MAX_CMD:-10}'
		set lucicodex.main.log_file='${OLD_LOG:-/tmp/lucicodex.log}'
		commit lucicodex
	EOF
	
	# Clean up old anonymous sections
	uci -q delete lucicodex.@api[0] 2>/dev/null
	uci -q delete lucicodex.@settings[0] 2>/dev/null
	uci -q commit lucicodex 2>/dev/null
fi

# Set proper permissions on config file (readable by uhttpd/luci)
chmod 644 /etc/config/lucicodex 2>/dev/null || true

exit 0
